# docker/systemd/taskrunner-vault.service
[Unit]
Description=TaskRunner.Vault Family
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/yj-family/TaskRunner.Vault
EnvironmentFile=/etc/systemd/system/taskrunner-vault.env
Environment=TASKRUNNER_DATA_DIR=/opt/yj-family/data
ExecStart=/usr/bin/dotnet TaskRunner.Vault.dll --urls http://0.0.0.0:8790
Restart=always
RestartSec=5
StandardOutput=append:/opt/yj-family/logs/taskrunner-vault.log
StandardError=append:/opt/yj-family/logs/taskrunner-vault.log

[Install]
WantedBy=multi-user.target
